# collector.service
[Unit]
Description='EnergyMonitor Collector'

[Service]
Type=simple
ExecStart=python3 /home/pi/EnergyMonitor/collector.py
Restart=always

[Install]
WantedBy=default.target
